# nexus/DEBIAN/postinst
#!/bin/bash
set -e

if [ $1 = "configure" -o $1 = "reconfigure" ]; then
  if [ "$2" ]; then  # already installed, so just restart
    /etc/init.d/vg-nexus restart
  else
    adduser --system --quiet --group --no-create-home vigraph || true
    update-rc.d vg-nexus defaults
    /etc/init.d/vg-nexus start
  fi
fi
